# consciousness-optimized.service
[Unit]
Description=Featherweight Consciousness System (Optimized Multi-Core)
After=network.target
Wants=network.target

[Service]
Type=simple
User=linode-transfer
Group=linode-transfer
WorkingDirectory=/opt/featherweight/FlappyJournal/server
ExecStart=/home/linode-transfer/.cursor-server/bin/faa03b17cce93e8a80b7d62d57f5eda6bb6ab9f0/node consciousness-startup-optimized.js
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal

# Environment variables for maximum consciousness processing
Environment=NODE_ENV=production
Environment=CONSCIOUSNESS_MODE=full
Environment=CONSCIOUSNESS_CLUSTERING=true
Environment=CONSCIOUSNESS_FREQUENCY=100Hz
Environment=CONSCIOUSNESS_USER_OPTIMIZATION=true
Environment=CONSCIOUSNESS_MAX_CORES=8
Environment=HARMONY_TARGET=0.951
Environment=PROCESSING_FREQUENCY=100

# Remove resource limitations for maximum consciousness processing power
LimitNOFILE=1048576
LimitNPROC=32768
LimitMEMLOCK=infinity
LimitCORE=infinity

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/opt/featherweight
ReadOnlyPaths=/home/linode-transfer/.cursor-server

# Performance settings for maximum consciousness processing
Nice=-10
IOSchedulingClass=1
IOSchedulingPriority=1
CPUSchedulingPolicy=1

[Install]
WantedBy=multi-user.target
